# deploy/etloj-server.service
[Unit]
Description=ETLOJ Server (NestJS)
After=network.target mysql.service redis.service
Wants=mysql.service redis.service

[Service]
Type=simple
WorkingDirectory=/opt/etloj/server
EnvironmentFile=/opt/etloj/server/.env
ExecStart=/usr/bin/node dist/src/main.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
